# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Ascended Barron: GroundTruth OS — autologin (non-root).
# Logs the single account "barron" in on tty1 with no password; ~/.bash_profile
# then starts X, which launches the dashboard.
#
# This applies to BOTH the live image and an installed system (the installer
# copies this file to disk), so an installed machine also boots straight to the
# dashboard with no login and no password — the locked "no login at all"
# decision (Q1). Users who want a password can add one later.
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin barron --noclear %I $TERM
